# cookbooks/nodejs-application/templates/default/supervisor.nodejs-application.conf.erb
#!upstart
description "<%= node[:app][:service][:name] %> server"
author      "Taliesin Sisson"

start on startup
stop on shutdown

script
  sudo touch /var/log/<%= node[:app][:service][:name] %>.log
  sudo chown <%= node[:app][:service][:user] %> /var/log/<%= node[:app][:service][:name] %>.log
  sudo chown -R <%= node[:app][:service][:user] %> <%= node[:app][:path] %>
  sudo -u <%= node[:app][:service][:user] %> <%= node[:nodejs][:root_dir] %>/bin/supervisor -x <%= node[:nodejs][:root_dir] %>/bin/node -w <%= node[:app][:path] %>/ <%= node[:app][:path] %>/server.js >> /var/log/<%= node[:app][:service][:name] %>.log 2>&1
end script
