# cookbooks/nodebox-development/templates/default/upstart.nodejs.conf.erb
#!upstart
description "<%= node[:app][:name] %> server"
author      "Taliesin Sisson"

start on startup
stop on shutdown

script
  sudo touch /var/log/<%= node[:app][:name] %>.log
  sudo chown <%= node[:nodejs][:service][:user] %> /var/log/<%= node[:app][:name] %>.log
  sudo -u <%= node[:nodejs][:service][:user] %> <%= node[:nodejs][:dir] %>/bin/supervisor -x <%= node[:nodejs][:dir] %>/bin/node -w <%= node[:app][:path] %>/ <%= node[:app][:path] %>/app.js >> /var/log/<%= node[:app][:name] %>.log 2>&1
end script
